# dot_config/systemd/user/mbsync.service.tmpl
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -c /home/{{.chezmoi.username}}/.config/isync/mbsyncrc -Va
ExecStartPost=/usr/bin/mu index
