# dot_config/systemd/user/goimapnotify@.service.d/override.conf
[Service]
Environment=NOTMUCH_CONFIG=%h/.config/notmuch/notmuchrc